# ocrquesseg.ini
[program:ocrquesseg]
command=gunicorn -c gunicorn.py manage:app
directory=/xdfapp/apps/ocrquesseg-svr/
user=root
autostart=true
autorestart=true
redirect_stderr=true
startsecs=20
startretries=3
stderr_logfile=/xdfapp/logs/ocrquesseg-svr/ocrquesseg-svr.err.log
stdout_logfile=/xdfapp/logs/ocrquesseg-svr/ocrquesseg-svr.out.log